# scripts/doomsday_option.service
[Unit]
Description=Doomsday Option Trading Service
After=network.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
User=trader
Group=trader
WorkingDirectory=/home/options_trading
Environment=PYTHONPATH=/home/options_trading
ExecStart=/home/options_trading/venv/bin/python -m scripts.main

# 基本限制
Restart=on-failure
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=30

# 资源限制
MemoryLimit=500M
TasksMax=100
LimitNOFILE=1024

# 基本安全限制
NoNewPrivileges=yes
PrivateTmp=yes
RestrictAddressFamilies=AF_INET AF_INET6

[Install]
WantedBy=multi-user.target